# profile.d/local07-mambarootprefix.sh
#!/bin/sh
MAMBA_ROOT_PREFIX="/opt/lsst/software/stack/conda"
export MAMBA_ROOT_PREFIX
