# static/etc/profile.d/local09-umask.sh
#!/bin/sh
if [ "${RSP_SITE_TYPE}" == "science" ]; then
    umask 002
fi
