# static/etc/profile.d/local07-setupstack.sh
#!/bin/sh
if [ -n "${RUNNING_INSIDE_JUPYTERLAB}" ]; then
    . /opt/lsst/software/stack/loadLSST.bash
fi
